# conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc --cert /etc/ssl/novnc/self-signed.pem 8080 localhost:5900
autorestart=true
