# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -passwdfile %(ENV_HOME)s/.vnc/passwd
autorestart=true
